# com.thisisgulshanshah.space-for-utm.autogitpush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thisisgulshanshah.space-for-utm.autogitpush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/thisisgulshanshah/Documents/SPACE FOR UTM/scripts/auto_git_push.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/thisisgulshanshah/Documents/SPACE FOR UTM</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/Users/thisisgulshanshah/Documents/SPACE FOR UTM/.auto-git-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/thisisgulshanshah/Documents/SPACE FOR UTM/.auto-git-launchd.err.log</string>
</dict>
</plist>
